# packaging/debian/postinst.in
#!/bin/sh
# postinst dla pakietu retractordb: tworzy użytkownika usługi i włącza jednostkę
# systemd (bez natychmiastowego startu — wystartuje przy następnym boocie).
#
# UWAGA: plik generowany z postinst.in przez configure_file() — @RETRACTOR_QUERY_FILE@
# podstawiane z CMake (jedno źródło prawdy, wspólne z ExecStart jednostki). Nie edytuj
# wygenerowanej kopii; zmiany nanoś w szablonie .in.
set -e

case "$1" in
  configure)
    # Użytkownik systemowy 'retractor' (User=/Group= w jednostce).
    if ! getent passwd retractor >/dev/null 2>&1; then
      if command -v adduser >/dev/null 2>&1; then
        adduser --system --group --no-create-home --home /nonexistent \
          --gecos "RetractorDB service" retractor || true
      else
        useradd --system --user-group --no-create-home \
          --shell /usr/sbin/nologin retractor || true
      fi
    fi

    # Przykład konfiguracji do /etc/retractor (sshd-style). Kopiujemy tylko gdy go
    # tam nie ma — nie nadpisujemy zmian administratora. Żywy retractor.toml NIE jest
    # tworzony: brak konfiguracji to poprawny stan (usługa rusza z domyślnymi).
    if [ -f /usr/share/retractordb/retractor.toml.example ]; then
      mkdir -p /etc/retractor
      if [ ! -f /etc/retractor/retractor.toml.example ]; then
        cp /usr/share/retractordb/retractor.toml.example /etc/retractor/retractor.toml.example
      fi
    fi

    # Przykład env-file jednostki (EnvironmentFile). Tylko .example — żywy service.env
    # NIE jest tworzony: jego brak to poprawny stan (jednostka używa domyślnej ścieżki).
    if [ -f /usr/share/retractordb/service.env.example ]; then
      mkdir -p /etc/retractor
      if [ ! -f /etc/retractor/service.env.example ]; then
        cp /usr/share/retractordb/service.env.example /etc/retractor/service.env.example
      fi
    fi

    # Kanoniczny plik zapytań wskazywany przez ExecStart jednostki. Pusty plik = tryb idle.
    # Tworzymy tylko gdy go nie ma (nie nadpisujemy zapytań administratora). Właściciel to
    # użytkownik usługi, by nowa instancja mogła go nadpisać przy przekazaniu zapytania (restart).
    mkdir -p "$(dirname @RETRACTOR_QUERY_FILE@)"
    if [ ! -f @RETRACTOR_QUERY_FILE@ ]; then
      : > @RETRACTOR_QUERY_FILE@
    fi
    chown retractor:retractor "$(dirname @RETRACTOR_QUERY_FILE@)" @RETRACTOR_QUERY_FILE@ 2>/dev/null || true

    # Rejestracja jednostki w systemd (tylko gdy systemd działa).
    if [ -d /run/systemd/system ]; then
      systemctl daemon-reload || true
      systemctl enable xretractor.service || true
    fi
    ;;
esac

exit 0
