# packaging/systemd/xretractor.service
# Przykładowa jednostka systemd dla serwera RetractorDB (xretractor).
#
# Instalacja:
#   sudo cp packaging/systemd/xretractor.service /etc/systemd/system/
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin retractor   # jeśli nie istnieje
#   sudo mkdir -p /etc/retractor && sudo chown retractor: /etc/retractor
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now xretractor
#
# Logi:   journalctl -u xretractor -f
# Status: systemctl status xretractor   (lub: xretractor --status)

[Unit]
Description=RetractorDB continuous query engine (xretractor)
Documentation=https://dokumentacja.retractordb.com
After=network.target

[Service]
Type=simple
User=retractor
Group=retractor

# Tryb usługowy: log na stderr (przechwytywany przez journald), bez TTY.
# Start bez pliku .rql uruchamia tryb bezczynny (idle); aby wczytać zapytania,
# wskaż plik jako pierwszy argument, np. ExecStart=... --service --noanykey /etc/retractor/startup.rql
#
# Tryb usługowy logowania można też włączyć zmienną środowiskową zamiast flagi --service:
#   Environment=XRETRACTOR_SERVICE=1
ExecStart=/usr/local/bin/xretractor --service --noanykey

KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
